# ops/launchd/com.healthos.runtime-async.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.healthos.runtime-async</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>/opt/healthos/ts/packages/runtime-async/dist/index.js</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/Shared/HealthOS/logs/runtime-async.out.log</string>
  <key>StandardErrorPath</key><string>/Users/Shared/HealthOS/logs/runtime-async.err.log</string>
</dict>
</plist>
